# d.conf
[program:d]
command=python3 /home/curator2/d/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/d/test.err.log
stdout_logfile=/home/curator2/d/test.out.log